# deploy/bird/bird.service.d/restart.conf
# Rotko: keep BIRD alive across crashes.
#
# 2026-07-16 incident: BIRD 3.2.1 segfaulted on bkk12 during a full-table
# reconvergence. bkk12's unit had Restart=no, so one SIGSEGV left the daemon
# dead and the validator's /32 unreachable for ~4h. The packaged unit elsewhere
# uses Restart=on-abort with StartLimitBurst=5/10s, which gives up after a
# crash-loop. Both are wrong for a routing daemon whose death blackholes a node.
#
# Install: place at /etc/systemd/system/bird.service.d/restart.conf on every
# BIRD host, then `systemctl daemon-reload` (does NOT restart the running bird).
[Unit]
# Never give up: disable the start-rate limiter so systemd always retries.
StartLimitIntervalSec=0

[Service]
# Restart on any unclean exit (segfault, non-zero, timeout), with a short backoff
# so a persistent bad condition recovers the instant it clears instead of thrashing.
Restart=on-failure
RestartSec=2s
